# project/supervisord.conf
[program:video_capturing]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py video_capturing